# conf/supervisor.d/apache.conf
[program:apache]
command=/usr/bin/apachectl -DFOREGROUND
