# conf/supervisor.d/rt-backup.conf
[program:rt-backup]
command=/usr/bin/bash -c "while true; do rsync -a /usr/share/webapps/rutorrent/share/settings/* /app/rutorrent/; sleep 15; done"
user=torrents
