# conf/supervisor.d/rtorrent.conf
[program:rtorrent]
command=/usr/bin/rtorrent -d /torrents -o scgi_port=localhost:5000
user=torrents
